# overlays/parallella-cfg/etc/rc.local
version https://git-lfs.github.com/spec/v1
oid sha256:aa2d01a3cc013433f0d165a9c3e62c5473168e85262214f3ab1cbcdc445cb78b
size 586
